# deploy/flowdmr-sidecar.service
[Unit]
Description=FlowDMR sidecar (DMR -> TETRA local injector)
Documentation=https://github.com/razvanzeces/flowdmr
# Start after FlowStation so the entity's IPC socket is already listening.
After=network-online.target tetra.service
Wants=network-online.target

[Service]
Type=simple
# Runs as root (like tetra.service) so the RTL-SDR USB device and dsd-neo just work.
# PATH must include /usr/local/bin so the child dsd-neo is found.
Environment=PATH=/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/flowdmr-sidecar --config /etc/flowdmr/sidecar.toml
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/flowdmr/flowdmr.log
StandardError=append:/var/log/flowdmr/flowdmr.log

[Install]
WantedBy=multi-user.target
