# deploy/flowdmr-sidecar.service
[Unit]
Description=FlowDMR sidecar (DMR -> TETRA local injector)
Documentation=https://github.com/razvanzeces/flowstation
# Start after FlowStation so the IPC entity socket is already listening.
After=network-online.target bluestation-bs.service
Wants=network-online.target

[Service]
Type=simple
# Adjust paths to where you installed the binary and config.
ExecStart=/usr/local/bin/flowdmr-sidecar --config /etc/flowdmr/sidecar.toml
Restart=on-failure
RestartSec=3
# dsd-neo needs USB (RTL-SDR) access; run as a user in the plugdev group.
User=flowdmr
Group=plugdev
Environment=RUST_LOG=info
# Localhost-only service; lock it down.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/flowdmr
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
